# configuração-do-etc/etc/profile.d/01-gcc.sh
# Configuração do GCC
export CC=gcc
export CXX=g++
export PATH="/usr/lib/gcc/bin:$PATH"
export MANPATH="/usr/share/gcc-data/man:$MANPATH"
