# configuração-do-etc/etc/profile.d/02-python.sh
# Configuração do Python
export PYTHONPATH="/usr/lib/python3.11/site-packages"
export PATH="/usr/lib/python3.11/bin:$PATH"
