# configuração-do-etc/etc/profile.d/10-aliases.sh
# Aliases globais
alias ll='ls -lh --color=auto'
alias grep='grep --color=auto'
